# backup/backup.timer
[Unit]
Description=Back up folders into Backup 

[Timer]
OnCalendar=weekly
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target
